# micro_httpd.socket
[Unit]
Description=Socket for micro_httpd

[Socket]
ListenStream=8912
Accept=yes

[Install]
WantedBy=sockets.target
